# yacs.cron
SHELL=/bin/bash
*/10 * * * * cd %(projectpath)s && %(python)s manage.py import_course_data 2>&1 >> %(logpath)s/import_course_data.log && %(python)s manage.py update_index --age 1 2>&1 >> %(logpath)s/elasticsearch_update.log
*/25 * * * * cd %(projectpath)s && %(python)s manage.py create_section_cache 2>&1 >> %(logpath)s/create_section_cache.log
